# efi/mkosi/fedora/arm64/mkosi.extra/etc/rc.d/rc.local
#!/bin/bash

set -x

# Start `systemd-networkd`
systemctl enable --now systemd-networkd.service
